# carekeeper-restart.service
[Unit]
Description=Restart CareKeeper service (nightly, to keep memory/cache fresh on Pi5)

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart carekeeper.service
